# alpine-based/zerotier/access/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:zerotier]
command=zerotier-one -U
autorestart=true
priority=1

[program:ssh]
command=/bin/sh /ssh.sh
redirect_stderr=true
priority=2
